# templates/pet-waker.service
[Unit]
Description=Pet Waker - Wake-on-request proxy for pet-cli
After=network.target

[Service]
Type=simple
WorkingDirectory={{PET_DIR}}/waker
ExecStart=/usr/bin/node {{PET_DIR}}/waker/pet-waker.js
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PET_DIR={{PET_DIR}}
Environment=PET_CONFIG_DIR={{PET_CONFIG_DIR}}
Environment=PET_WAKER_PORT=3999
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pet-waker

[Install]
WantedBy=default.target
